# fisheye_webservice/util/upstart/fisheye_webserver.conf
description "Fisheye webservice encapsulated into docker container"
author      "Dmitry Bielous<dmigous@gmail.com>"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  FISHEYE_REPO_PATH=/home/dmigous/workspace/fisheye/fisheye_webservice/util/upstart/../../../ # this path auto substituted from install_upstart_script.sh
  HAS_STOPPED_FEDORA_FISHEYE_CONTAINER=`docker ps -a | grep "dmigous/fedora_fisheye" -m1 | awk '{print $1}' | wc -l`
  if HAS_STOPPED_FEDORA_FISHEYE_CONTAINER; then
    /usr/bin/docker start -a fisheye_webserver
  else
    docker run --name fisheye_webserver -i -t -p 80:80 -v \
      $FISHEYE_REPO_PATH:/mnt dmigous/fedora_fisheye /bin/bash \
      -c "source /root/.bashrc; /mnt/fisheye_webservice/util/prepare_apache.sh; /etc/httpd/run_apache_foreground"
  fi
end script
